# src/etc/s6/services/cron/finish
#!/usr/bin/execlineb -S2

s6-echo "[finish] cron: exit=${1} signal=${2}"
