# src/fs/etc/services.d/rsyslog_tail/finish
#!/usr/bin/execlineb -S2

s6-echo "[finish] rsyslog_tail: exit=${1} signal=${2}"
