# src/etc/s6/services/rsyslog_tail/run
#!/usr/bin/execlineb -P

trap
{
  term
  {
    foreground {
        loopwhilex -x 0
        foreground { s6-sleep 1 }
        pipeline { tail -n 100 /var/log/syslog } redirfd -w 1 /dev/null grep "rsyslogd:.*exiting on"
    }
    importas -u ! ! kill ${!}
  }
}
tail -F -n 0 /var/log/syslog
