# src/etc/s6/services/rsyslog_tail/run
#!/usr/bin/execlineb -P

trap
{
  term
  {
    foreground { s6-svwait -D /var/run/s6/services/rsyslog }
    importas -u ! ! kill ${!}
  }
}
tail -F -n 0 /var/log/syslog
