# src/etc/s6/services/rsyslog/run
#!/bin/sh

rm -f /var/run/rsyslogd.pid
exec /usr/sbin/rsyslogd -n -i /var/run/rsyslogd.pid
