# src/etc/s6/services/rsyslog/run
#!/bin/sh

exec /usr/sbin/rsyslogd -n -i /var/run/rsyslogd.pid
